# data/supervisor.conf
[supervisord]
nodaemon=true

[program:ser2net]
directory=/usr/sbin/
command=/bin/bash -c "ser2net -d -u"

[program:klippy]
directory=/home/klippy
command=/bin/bash -c "./klippy-env/bin/python ./klipper/klippy/klippy.py ./.config/printer.cfg"

#CMD ["/bin/sh", "/usr/sbin/ser2net -d -u && /home/klippy/klippy-env/bin/python /home/klippy/klipper/klippy/klippy.py /home/klippy/.config/printer.cfg"]